# net/sshguard/files/sshguard.init
#!/bin/sh /etc/rc.common
# Copyright (C) 2017 OpenWrt.org

USE_PROCD=1

START=90
STOP=10

PROG=/usr/sbin/sshguard
PIDFILE="`uci -q get sshguard.settings.pid_file`"
PIDFILE="${PIDFILE:-/var/run/sshguard.pid}"

start_service() {
	procd_open_instance
	procd_set_param command $PROG -i $PIDFILE
	procd_set_param pidfile $PIDFILE
	procd_close_instance
}
